# override-files/etc/uci-defaults/200-disable-firmware-updates
#!/bin/sh

uci set glconfig.autoupdate.firmware_path='http://www.google.com'
uci commit
